# pi-home/.config/autostart/autovlc.desktop
[Desktop Entry] 
Type=Application
Exec=clementine /home/pi/Desktop/SomaFM-Radio-Stations.xspf
Name[en_GB]=autovlc.desktop

