# pi-home/.config/autostart/auto-start-aconngui.desktop
[Desktop Entry] 
Type=Application
Exec=aconnectgui
